# engine/birdengine-recording.service
[Unit]
Description=BirdEngine Recording — Audio capture
After=sound.target

[Service]
Type=simple
User=BIRDASH_USER
ExecStart=BIRDASH_HOME/birdash/engine/record.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
