# engine/ttyd.service
[Unit]
Description=Birdash Web Terminal (ttyd)
After=network.target

[Service]
Type=simple
User=BIRDASH_USER
ExecStart=/usr/local/bin/ttyd --port 7681 --base-path /birds/terminal --writable --interface lo bash
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
