# config/birdash.service
[Unit]
Description=Birdash API Server
After=network.target

[Service]
Type=simple
User=BIRDASH_USER
WorkingDirectory=BIRDASH_HOME/birdash
ExecStart=/usr/bin/node BIRDASH_HOME/birdash/server/server.js
Restart=on-failure
RestartSec=5
Environment=BIRDASH_PORT=7474
Environment=BIRDASH_DB=BIRDASH_HOME/birdash/data/birds.db
# Environment=EBIRD_API_KEY=your_key_here
# Environment=BW_STATION_ID=your_station_id

[Install]
WantedBy=multi-user.target
